# run-after-restore.service
[Unit]
Description=Run bootstrap script after system restore
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/declare-sh/bootstrap.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
